# get_pw.cron
*/30 * * * * env - `cat /etc/environment` node /base/getPw.js >> /dev/null 2>&1

